# net-proxy/ssl-proxy/files/ssl-proxy.socket
[Unit]
Description=Fuck gfw on squid with ssl
Conflicts=ssl-proxy.service

[Socket]
ListenStream=8080
Accept=yes

